# mOS/tools/yod/02-mos.rules
ENV{MOS_SYSFS_FIXED}=="", RUN+="/bin/chmod -R o=g /sys/kernel/mOS", ENV{MOS_SYSFS_FIXED}="yes"
